# restartnfc.service
[Unit]
Description=Restart ReadNFC

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart readnfc.service

